# bin/service/nginx/run
#!/bin/sh
exec 2>&1
exec service nginx start 2>&1 > /dev/null
